# GNULinux/config/autostart/syncthing.desktop
[Desktop Entry]
Name=Syncthing
Exec=syncthing --no-browser
Type=Application
