# GNULinux/config/autostart/upgrade_all.desktop
[Desktop Entry]
Name=Upgrade All
Exec=/home/summonholmes/.local/bin/upgrade_all.sh
Type=Application
